# operating-system/orderpy.service
[Unit]
Description=Daemon for order Python service
After=network.target

[Service]
Restart=always
WorkingDirectory=/home/mriansyah93/python-projects/uho-parallel-python
Environment="PYTHONPATH=/home/mriansyah93/.local/lib/python3.8/site-packages"
ExecStart=/home/mriansyah93/.local/bin/uvicorn main:app --reload --port 7080

[Install]
WantedBy=multi-user.target